# supervisord.conf
[supervisord]
nodaemon=true

[program:scrapyd]
command=scrapyd

[program:logparser]
command=logparser -dir /scrapyd/logs
